# rootfs/etc/s6-overlay/s6-rc.d/tvheadend/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: TVHeadend
# Runs TVHeadend script
# ==============================================================================

bashio::log.info "Starting TVHeadend..."

exec /usr/bin/run_tvheadend.sh
